# deploy/titan-dashboard.service
# ============================================================
#  Titan Dashboard — systemd service
#  Runs Flask dashboard on :5000, accessible at http://IP:5000
# ============================================================
[Unit]
Description=Titan Trading Dashboard (Flask)
After=network.target titan-bot.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/titan/agentic_trader
Environment=PATH=/home/ubuntu/titan/venv/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/ubuntu/titan/agentic_trader/.env

# Gunicorn with gthread worker for SSE streaming support
ExecStart=/home/ubuntu/titan/venv/bin/gunicorn \
    --bind 0.0.0.0:5000 \
    --workers 2 \
    --threads 4 \
    --worker-class gthread \
    --timeout 120 \
    dashboard:app

Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=300

StandardOutput=append:/home/ubuntu/titan/logs/dashboard.log
StandardError=append:/home/ubuntu/titan/logs/dashboard-error.log

[Install]
WantedBy=multi-user.target
