# agentic_trader/deploy/titan-start.service
[Unit]
Description=Start Titan trading bot (triggered by timer)

[Service]
Type=oneshot
ExecStart=/bin/systemctl start titan
