# agentic_trader/deploy/titan-stop.service
[Unit]
Description=Stop Titan trading bot (triggered by timer)

[Service]
Type=oneshot
ExecStart=/bin/systemctl stop titan
